# ansible/roles/docker/files/docker.service
.include /usr/lib/systemd/system/docker.service
[Service]
MountFlags=shared
EOF